# cloud-node/services/scheduled-dispatcher.timer
[Unit]
Description=Run scheduled task dispatcher every minute

[Timer]
Unit=scheduled-dispatcher.service
OnCalendar=*-*-* *:*:00
Persistent=true
AccuracySec=5s

[Install]
WantedBy=timers.target
