# cloud-node/services/scheduled-dispatcher.service
[Unit]
Description=Browning Cloud Node — Scheduled Task Dispatcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=agent
Group=agent
ExecStart=/home/agent/claude-mem/cloud-node/scripts/scheduled-dispatcher.sh
EnvironmentFile=/home/agent/.config/task-watcher.env
TimeoutStartSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=scheduled-dispatcher
